# supervisord.conf
[program:daemonname]
command = python coe.py
directory = /home/dotcloud/current/
TZ = Asia/Singapore