# supervisord.conf
[program:daemonname]
command = python ~/current/coe.py -d
directory = /home/dotcloud/current/
TZ = Asia/Singapore